# deploy/systemd/ai-analysis-acceptance-healthcheck.service
[Unit]
Description=Library AI Analysis trace acceptance healthcheck
After=network-online.target library.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/library/app
EnvironmentFile=-/etc/default/library-monitoring
ExecStart=/usr/bin/env npm run acceptance:healthcheck -- --json
