# deploy/systemd/ai-analysis-ui-smoke-healthcheck.service
[Unit]
Description=Library AI Analysis browser smoke healthcheck
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/library/app
EnvironmentFile=-/etc/default/library-monitoring
Environment=NODE_ENV=production
ExecStart=/usr/bin/env npm run ui:smoke:healthcheck -- --json
