# deploy/systemd/ai-analysis-ui-qa-healthcheck.timer
[Unit]
Description=Periodic Library AI Analysis authenticated browser QA healthcheck

[Timer]
OnBootSec=6min
OnUnitActiveSec=6h
Unit=ai-analysis-ui-qa-healthcheck.service
Persistent=true

[Install]
WantedBy=timers.target
