# deploy/systemd/library-system-healthcheck.service
[Unit]
Description=Library service-chain healthcheck
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/library/app
EnvironmentFile=-/etc/default/library-monitoring
ExecStart=/usr/bin/env npm run healthcheck -- --json
